# dotfiles/.config/systemd/user/switch-to-light.service
[Unit]
Description=Switch to light theme

[Service]
Type=oneshot
ExecStart=/home/<user>/arch-rice/scripts/switch-theme.sh light-theme

